# cron/crontab
MAILTO=webmaster@digitalkingdom.org
37 */5 * * * /home/spdkhaproxy/dk-haproxy/rebuild_images.sh
1  1   * * 4 /home/spdkhaproxy/dk-haproxy/renew_certs.sh
